# scripts/com.claude-figma-mcp.socket.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.claude-figma-mcp.socket</string>

    <key>ProgramArguments</key>
    <array>
        <string>$PROJECT_PATH/scripts/Claude Figma MCP.app/Contents/MacOS/Claude Figma MCP</string>
    </array>

    <key>WorkingDirectory</key>
    <string>$PROJECT_PATH</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>$HOME/Library/Logs/claude-figma-mcp-socket.log</string>

    <key>StandardErrorPath</key>
    <string>$HOME/Library/Logs/claude-figma-mcp-socket.error.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>BUN_EXEC_PATH</key>
        <string>$BUN_PATH</string>
        <key>BUN_PROJECT_PATH</key>
        <string>$PROJECT_PATH</string>
    </dict>
</dict>
</plist>
